# .local/bin/polybar/95-battery.rules
SUBSYSTEM=="power_supply", ENV{POWER_SUPPLY_ONLINE}=="0", \
    RUN+="/home/z/.local/bin/polybar/poly-bat --update"
SUBSYSTEM=="power_supply", ENV{POWER_SUPPLY_ONLINE}=="1", \
    RUN+="/home/z/.local/bin/polybar/poly-bat --update"
